# config/systemd/user/checkupdates.service
[Unit]
Description=Check for pacman updates

[Service]
Type=oneshot
ExecStart=/home/hexid/.local/bin/updates.sh
